# polymarket-bot.service
[Unit]
Description=Polymarket Market Making Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polymarket
WorkingDirectory=/opt/polymarket-bot
ExecStart=/opt/polymarket-bot/venv/bin/python main.py

# Auto-restart on crash (5 second cooldown)
Restart=on-failure
RestartSec=5

# Kill cleanly — SIGINT first (triggers graceful shutdown + order cancellation),
# then SIGKILL after 30s if the bot doesn't exit.
KillSignal=SIGINT
TimeoutStopSec=30

# Environment
EnvironmentFile=/opt/polymarket-bot/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/polymarket-bot/logs /opt/polymarket-bot/positions.json
ProtectHome=true
PrivateTmp=true

# Resource limits — prevent runaway memory from crashing the server
MemoryMax=512M
MemoryHigh=384M

[Install]
WantedBy=multi-user.target
